# docker/supervisor/conf.d/cron.conf
[program:cron]
command=cron && tail -f /var/log/cron.log
autostart=true
autorestart=true
redirect_stderr=false
